# docker/rootfs/etc/supervisord.conf
; Dockery process supervision.
;
; Three long-running processes managed by supervisord (PID 1):
;   - dockery-api (priority 10): boots first — opens SQLite, generates
;     the Ed25519 keypair at /data/config/jwt-*.pem, listens on :3001.
;   - registry   (priority 20): blocks until jwt-public.pem exists on
;     disk, then starts. This wait-for-file is Dockery's substitute for
;     s6-overlay's notification-fd; when M2.4 flips registry auth to
;     token mode, the public key MUST be present or registry refuses
;     to start.
;   - nginx      (priority 30): fronts all three on :5000. Safe to
;     start even if upstreams are briefly unreachable (it will just
;     return 502 until they listen).
;
; All services run as root inside an isolated container; logs go
; straight to the container's stdout/stderr for `docker logs`.

[supervisord]
nodaemon=true
user=root
logfile=/dev/null
logfile_maxbytes=0
pidfile=/run/supervisord.pid

[unix_http_server]
file=/run/supervisor.sock
chmod=0700

[supervisorctl]
serverurl=unix:///run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

; ============== dockery-api ==============
; Wrap with sh -c so we can mkdir the subdirs of /data on first launch —
; /data is a docker volume and therefore empty before first boot.
; biz/keystore.go creates /data/config on demand, but SQLite won't
; create /data/db for us (error 14 "unable to open database file"),
; so we do it explicitly here.
[program:dockery-api]
command=/bin/sh -c 'mkdir -p /data/db /data/config /data/registry && exec /usr/local/bin/dockery-api -conf /etc/dockery'
priority=10
autostart=true
autorestart=true
startsecs=1
stopsignal=TERM
stopwaitsecs=10
; Ensure child-groups (any accidental forks) die together on stop.
stopasgroup=true
killasgroup=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

; ============== registry ==============
; Wait-for-file substitute for s6-overlay's notification-fd: block
; until dockery-api has written the public key, then exec registry.
; The loop polls every 200ms and gives up after ~30s to avoid masking
; a genuinely broken api boot.
[program:registry]
command=/bin/sh -c '\
  i=0; \
  until [ -f /data/config/jwt-jwks.json ]; do \
    i=$((i+1)); \
    [ "$i" -gt 150 ] && echo "timed out waiting for jwt-jwks.json" >&2 && exit 1; \
    sleep 0.2; \
  done; \
  exec /usr/local/bin/registry serve /etc/docker/registry/config.yml'
priority=20
autostart=true
autorestart=true
startsecs=2
stopsignal=TERM
stopwaitsecs=15
stopasgroup=true
killasgroup=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

; ============== nginx ==============
[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
priority=30
autostart=true
autorestart=true
startsecs=1
stopsignal=QUIT
stopwaitsecs=10
stopasgroup=true
killasgroup=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
